# cmd/clex/templates/com.reissui.clexd.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd agent for the clex daemon (macOS). Generated by `clex service install`.
  The daemon needs no open ports (it long-polls GitHub and Telegram). KeepAlive
  restarts it on unclean exit; clex is crash-safe and reconstructs in-flight work
  from GitHub labels on startup. On stop, launchd sends SIGTERM for graceful
  shutdown (cancel runners, revert labels, flush the store).
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>{{.Label}}</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{.BinaryPath}}</string>
        <string>--repo</string>
        <string>{{.Repo}}</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ExitTimeOut</key>
    <integer>30</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>{{.LogPath}}</string>
    <key>StandardErrorPath</key>
    <string>{{.ErrLogPath}}</string>
</dict>
</plist>
