# go/cmd/clex/templates/clexd.service.tmpl
# systemd unit for the clex daemon (Linux). Generated by `clex service install`.
#
# The daemon needs no open ports (it long-polls GitHub and Telegram). Provide
# credentials via the EnvironmentFile below so the token never lives in this
# unit; that file must be 0600 and owned by the clex user.

[Unit]
Description=clex daemon (agentic development orchestrator)
Documentation=https://github.com/reissui/clex
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{.User}}
Group={{.User}}
EnvironmentFile={{.EnvFile}}
ExecStart={{.BinaryPath}} --repo {{.Repo}}
Restart=on-failure
RestartSec=5s
KillSignal=SIGTERM
TimeoutStopSec=30s

# --- Hardening (least privilege) ---
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{.Home}}
PrivateTmp=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
